# stock-service.service
[Unit]
Description=Stock Screening System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pz
Environment=PATH=/root/pz/.venv/bin
ExecStart=/root/pz/.venv/bin/python /root/pz/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stock-service

[Install]
WantedBy=multi-user.target
